# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/supervisord.log
pidfile=/var/run/supervisord.pid

[program:web]
command=python app.py
directory=/app
stdout_logfile=/app/data/web.log
stderr_logfile=/app/data/web_error.log
autostart=true
autorestart=true
startretries=3
priority=100
environment=PORT=%(ENV_PORT)s

[program:daemon]
command=python daemon.py
directory=/app
stdout_logfile=/app/data/daemon.log
stderr_logfile=/app/data/daemon_error.log
autostart=true
autorestart=true
startretries=3
priority=200
